# NAT Instance/libs/nat.service
[Unit]
Description=NAT Service for AWS
 
[Service]
Type=simple
ExecStart=/usr/local/bin/nat-setup start
ExecStop=/usr/local/bin/nat-setup stop
RemainAfterExit=yes
Restart=no

[Install]
WantedBy=network-up.target